# systemd/user/altcoin-trend.service
[Unit]
Description=Altcoin Trend System Daemon
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/altcoin-trend-system
EnvironmentFile=-%h/.config/acts/acts.env
Environment=PYTHONPATH=%h/altcoin-trend-system/src
ExecStart=%h/altcoin-trend-system/.venv/bin/python -m altcoin_trend.daemon
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
